# nas_ports/sysutils/syslog-ng/files/syslog-ng.in
#!/bin/sh
#
# $FreeBSD: head/sysutils/syslog-ng/files/syslog-ng.in 340872 2014-01-24 00:14:07Z mat $
#

# PROVIDE: syslogd
# REQUIRE: mountcritremote cleanvar ldconfig
# BEFORE:  SERVERS

#
# Add the following line to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
#syslog_ng_enable="YES"
#syslog_ng_config="-u daemon"
#syslog_ng_pid="/var/run/syslog-ng.pid"
#

syslog_ng_purgeklog=${syslog_ng_purgeklog-"NO"}

. /etc/rc.subr

name=syslog_ng
rcvar=syslog_ng_enable

command=%%PREFIX%%/sbin/syslog-ng
required_files=%%PREFIX%%/etc/syslog-ng.conf
extra_commands=reload

if checkyesno syslog_ng_purgeklog; then
	start_precmd="echo \"Purging klog(9)\" && sysctl -w kern.msgbuf_clear=1 $start_precmd";
fi

stop_postcmd=stop_postcmd

load_rc_config $name

: ${syslog_ng_enable:="NO"}
pidfile=${syslog_ng_pid:-"/var/run/syslog.pid"}
command_args="-p ${pidfile} ${syslog_ng_config}"

stop_postcmd() {
  rm -f $pidfile
}

run_rc_command "$1"
