# src/freenas/etc/rc.d/ix-devd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-devd
# REQUIRE: ix-zfs
# BEFORE: devd pf

. /etc/rc.subr

#
# Generate CARP hooks right before mountlate.
#
ix_devd_gen()
{
	local IFS=\|

	local failovers=`${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} "SELECT COUNT(id) FROM system_failover;"`

	if [ ${failovers} -gt 0 ]; then
		cp /conf/base/etc/devd.conf /etc/devd.conf
		cat >> /etc/devd.conf << E*O*F

notify 100 {
   match "system"   "IFNET";
   match "subsystem" "carp.*";
   action "/bin/sh -x /etc/carp-state-change-hook \$subsystem \$type | logger";
};
E*O*F
		volumename=`${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} "SELECT vol_name FROM storage_volume WHERE vol_fstype = 'ZFS' AND id in (SELECT volume_id from system_failover ORDER BY -id LIMIT 1);"`

		# Failover script
		cat > /etc/carp-state-change-hook << E*O*F
#!/bin/sh

SUBSYSTEM=\$1
EVENT=\$2

if [ \${SUBSYSTEM} != "carp0" ]; then
	exit
fi

if [ \${EVENT} = "LINK_UP" ]; then
    echo "Entering UP on \${SUBSYSTEM}!" | logger

    touch -t \`date -v -1M +%Y%m%d%H%M.%S\` /tmp/heartbeat_barrier
    WASCONNECTED=\`find /tmp/heartbeat_state -mnewer /tmp/heartbeat_barrier | wc -l\`

    sleep 3
    STATUS0=\`ifconfig \${SUBSYSTEM} | grep "carp:" | awk '{print \$2}'\`
# Test if the heartbeat cable thinks we are master
    STATUS1=\`(ifconfig carp1 | grep carp: | awk '{print \$2;}' ; ifconfig carp2 | grep carp: | awk '{print \$2;}')| grep MASTER | wc -l\`
    STATUS2=\`(ifconfig carp1 | grep carp: | awk '{print \$2;}' ; ifconfig carp2 | grep carp: | awk '{print \$2;}')| grep BACKUP | wc -l\`

    echo \${STATUS0}:\${STATUS1}:\${STATUS2} | logger
    if [ \$STATUS0 != "MASTER" ]; then
	echo "Promoted then demoted, quitting." | logger
	# Just in case.  Demote ourselves.
	ifconfig \${SUBSYSTEM} advskew 100
        exit
    fi

    if [ \$STATUS1 -ne 2 -o \$STATUS2 -ne 0 ]; then
	if [ ! -e /tmp/coldboot ]; then
		echo "No coldboot seal!" | logger
		# It's Okay if heartbeat say we are the master.
		if [ \$STATUS1 -eq 1 -a \$STATUS2 -eq 1 ]; then
			STATUS_CARP1=\`ifconfig carp1 | grep "carp:" | awk '{print \$2}'\`
			if [ \$STATUS_CARP1 = "MASTER" ]; then
				# Consider ourselves the master in this case.
				# This assignment is not really needed, but
				# we set it here in case the daemon didn't do
				# it.
				WASCONNECTED=1
				echo "Elected master based on carp1/2 states" | logger
			else
				echo "Defeated as backup on carp1/2 states" | logger
				# Just in case.  Demote ourselves.
				ifconfig \${SUBSYSTEM} advskew 200
				exit
			fi
		else
			echo "WARNING! May be splitting brain!" | logger
			# Just in case.  Demote ourselves.
			ifconfig \${SUBSYSTEM} advskew 100
			exit
		fi
	else
		echo "Found coldboot seal, continuing!" | logger
		WASCONNECTED=1
	fi
    fi

    # Remove the seal here, to be safe.
    rm -f /tmp/coldboot;

    if [ \$WASCONNECTED -ne 1 ]; then
	exit
    fi

    ifconfig \${SUBSYSTEM} advskew 1
    /sbin/zpool import -o cachefile=none -R /mnt -f ${volumename}
    /sbin/zpool set cachefile=/data/zfs/zpool.cache ${volumename}

    #mount -uw /
    /etc/rc.d/statd restart
    /usr/sbin/service ix-nfsd quietstart
    /usr/sbin/service mountd restart
    /usr/sbin/service nfsd restart
    /usr/sbin/service ix-istgt quietstart
    /usr/sbin/service istgt forcestop
    /usr/sbin/service istgt restart
#
# There appears to be a small lag if we allow NFS traffic right away.  During
# this time, we fail NFS requests with ESTALE to the remote system.  This
# gives remote clients heartburn, so rather than try to deal with the
# downstream effect of that, instead we take a chill pill for 2 seconds.
#
    sleep 1
    echo "pass quick on \${SUBSYSTEM} keep state" > /etc/pf.conf
    pfctl -f /etc/pf.conf
    echo "$(date), $(hostname), assume master" | mail -s "Failover" root
fi

if [ \${EVENT} = "LINK_DOWN" ]; then
    echo "Entering DOWN!" | logger
    ifconfig \${SUBSYSTEM} advskew 100
    echo "set block-policy drop" > /etc/pf.conf.block
    echo "block drop in quick proto tcp from any to any port 860" >> /etc/pf.conf.block
    echo "block drop in quick proto udp from any to any port 860" >> /etc/pf.conf.block
    echo "block drop in quick proto tcp from any to any port 2049" >> /etc/pf.conf.block
    echo "block drop in quick proto udp from any to any port 2049" >> /etc/pf.conf.block
    echo "block drop in quick proto tcp from any to any port 3260" >> /etc/pf.conf.block
    echo "block drop in quick proto udp from any to any port 3260" >> /etc/pf.conf.block
    pfctl -f /etc/pf.conf.block
    /etc/rc.d/statd stop
    echo "Exporting ${volumename}!" | logger
    zpool export ${volumename} || watchdog -t 1
    echo "Exported ${volumename}!" | logger
    echo "$(date), $(hostname), assume backup" | mail -s "Failover" root
fi
E*O*F
		touch /tmp/heartbeat_state /tmp/coldboot
		chown nobody /tmp/heartbeat_state
	else
		/usr/bin/truncate -s 0 /etc/pf.conf
	fi
}

name="ix-devd"
start_cmd='ix_devd_gen'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
