# src/freenas/etc/systemd/systemd-udevd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/lib/systemd/systemd-udevd --event-timeout=300