# nas_ports/misc/freenas-files/pkg-plist
usr/pbi/.pbiwrapper-%%MACHINE_ARCH%%
usr/local/bin/extract-tarball
usr/local/bin/winacl
usr/local/bin/autotune
usr/local/sbin/pbi_create
usr/local/sbin/pbi_add
usr/local/sbin/pbi_addrepo
usr/local/sbin/pbi_browser
usr/local/sbin/pbi_autobuild
usr/local/sbin/pbi_delete
usr/local/sbin/pbi_deleterepo
usr/local/sbin/pbi_icon
usr/local/sbin/pbi_info
usr/local/sbin/pbi_indextool
usr/local/sbin/pbi_listrepo
usr/local/sbin/pbi_makepatch
usr/local/sbin/pbi_makeport
usr/local/sbin/pbi_makerepo
usr/local/sbin/pbi_metatool
usr/local/sbin/pbi_patch
usr/local/sbin/pbi_update
usr/local/sbin/pbi_update_hashdir
usr/local/sbin/pbid
usr/local/sbin/pbi-crashhandler
usr/local/share/pcbsd/distfiles/pcbsd.rpo
usr/local/bin/pc-extractoverlay
usr/local/share/pcbsd/distfiles/server-overlay.txz
usr/local/bin/warden
usr/local/share/warden/pluginjail-packages
usr/local/share/warden/pcbsd-utils-packages
usr/local/etc/warden.conf
usr/local/share/warden/debian-6-squeeze
usr/local/share/warden/gentoo-stage3-i486
usr/local/share/pcbsd/scripts/functions.sh
usr/local/share/warden/scripts/backend/checkstatus.sh
usr/local/share/warden/scripts/backend/checkupdates.sh
usr/local/share/warden/scripts/backend/chrootjail.sh
usr/local/share/warden/scripts/backend/createjail.sh
usr/local/share/warden/scripts/backend/createtemplate.sh
usr/local/share/warden/scripts/backend/cronsnap.sh
usr/local/share/warden/scripts/backend/deletejail.sh
usr/local/share/warden/scripts/backend/detail.sh
usr/local/share/warden/scripts/backend/exportjail.sh
usr/local/share/warden/scripts/backend/functions.sh
usr/local/share/warden/scripts/backend/importjail.sh
usr/local/share/warden/scripts/backend/jailid.sh
usr/local/share/warden/scripts/backend/listjails.sh
usr/local/share/warden/scripts/backend/listpbis.sh
usr/local/share/warden/scripts/backend/listpkgs.sh
usr/local/share/warden/scripts/backend/setupusers.sh
usr/local/share/warden/scripts/backend/startjail.sh
usr/local/share/warden/scripts/backend/stopjail.sh
usr/local/share/warden/scripts/backend/zfsclonesnap.sh
usr/local/share/warden/scripts/backend/zfscronsnap.sh
usr/local/share/warden/scripts/backend/zfslistclone.sh
usr/local/share/warden/scripts/backend/zfslistsnap.sh
usr/local/share/warden/scripts/backend/zfsmksnap.sh
usr/local/share/warden/scripts/backend/zfsrevertsnap.sh
usr/local/share/warden/scripts/backend/zfsrmclone.sh
usr/local/share/warden/scripts/backend/zfsrmsnap.sh
usr/local/share/warden/scripts/rc.d/shutdown.sh
usr/local/share/warden/scripts/rc.d/startup.sh
usr/local/share/warden/scripts/rc.d/wardenrc
usr/local/etc/wardenrc
usr/local/share/warden/agent.png
usr/local/share/warden/warden.png
usr/local/share/pcbsd/conf/iso-mirror
usr/local/share/pcbsd/conf/pbi-mirror
usr/local/share/pcbsd/conf/pcbsd-mirrors
usr/local/share/pcbsd/conf/pkg-mirror
usr/local/share/pcbsd/conf/update-mirror
