# nas_ports/sysutils/syslog-ng/pkg-plist
@unexec %D/etc/rc.d/syslog-ng.sh stop > /dev/null 2>&1 || true
@sample etc/scl.conf.sample
etc/syslog-ng.conf.dist
@sample etc/syslog-ng.conf.sample
bin/loggen
bin/pdbtool
bin/update-patterndb
include/syslog-ng/ack_tracker.h
include/syslog-ng/afinter.h
include/syslog-ng/alarms.h
include/syslog-ng/apphook.h
include/syslog-ng/atomic.h
include/syslog-ng/block-ref-parser.h
include/syslog-ng/bookmark.h
include/syslog-ng/cache.h
include/syslog-ng/cfg-args.h
include/syslog-ng/cfg-grammar.h
include/syslog-ng/cfg-lexer-subst.h
include/syslog-ng/cfg-lexer.h
include/syslog-ng/cfg-parser.h
include/syslog-ng/cfg-tree.h
include/syslog-ng/cfg.h
include/syslog-ng/children.h
include/syslog-ng/compat/compat.h
include/syslog-ng/compat/getutent.h
include/syslog-ng/compat/glib.h
include/syslog-ng/compat/lfs.h
include/syslog-ng/compat/pio.h
include/syslog-ng/compat/time.h
include/syslog-ng/compat/socket.h
include/syslog-ng/compat/string.h
include/syslog-ng/control/control-server.h
include/syslog-ng/control/control.h
include/syslog-ng/crypto.h
include/syslog-ng/dnscache.h
include/syslog-ng/driver.h
include/syslog-ng/file-perms.h
include/syslog-ng/filter/filter-call.h
include/syslog-ng/filter/filter-cmp.h
include/syslog-ng/filter/filter-expr-parser.h
include/syslog-ng/filter/filter-expr.h
include/syslog-ng/filter/filter-in-list.h
include/syslog-ng/filter/filter-netmask.h
include/syslog-ng/filter/filter-op.h
include/syslog-ng/filter/filter-pipe.h
include/syslog-ng/filter/filter-pri.h
include/syslog-ng/filter/filter-re.h
include/syslog-ng/filter/filter-tags.h
include/syslog-ng/gprocess.h
include/syslog-ng/gsockaddr.h
include/syslog-ng/gsocket.h
include/syslog-ng/host-resolve.h
include/syslog-ng/hostname.h
include/syslog-ng/iv.h
include/syslog-ng/iv_avl.h
include/syslog-ng/iv_event.h
include/syslog-ng/iv_event_raw.h
include/syslog-ng/iv_fd_pump.h
include/syslog-ng/iv_inotify.h
include/syslog-ng/iv_list.h
include/syslog-ng/iv_popen.h
include/syslog-ng/iv_signal.h
include/syslog-ng/iv_thread.h
include/syslog-ng/iv_tls.h
include/syslog-ng/iv_wait.h
include/syslog-ng/iv_work.h
include/syslog-ng/ivykis/iv.h
include/syslog-ng/ivykis/iv_avl.h
include/syslog-ng/ivykis/iv_event.h
include/syslog-ng/ivykis/iv_event_raw.h
include/syslog-ng/ivykis/iv_fd_pump.h
include/syslog-ng/ivykis/iv_list.h
include/syslog-ng/ivykis/iv_popen.h
include/syslog-ng/ivykis/iv_signal.h
include/syslog-ng/ivykis/iv_thread.h
include/syslog-ng/ivykis/iv_tls.h
include/syslog-ng/ivykis/iv_wait.h
include/syslog-ng/ivykis/iv_work.h
include/syslog-ng/libtest/mock-transport.h
include/syslog-ng/libtest/msg_parse_lib.h
include/syslog-ng/libtest/persist_lib.h
include/syslog-ng/libtest/proto_lib.h
include/syslog-ng/libtest/template_lib.h
include/syslog-ng/libtest/testutils.h
include/syslog-ng/logmatcher.h
include/syslog-ng/logmpx.h
include/syslog-ng/logmsg.h
include/syslog-ng/logpipe.h
include/syslog-ng/logproto/logproto-buffered-server.h
include/syslog-ng/logproto/logproto-builtins.h
include/syslog-ng/logproto/logproto-client.h
include/syslog-ng/logproto/logproto-dgram-server.h
include/syslog-ng/logproto/logproto-framed-client.h
include/syslog-ng/logproto/logproto-framed-server.h
include/syslog-ng/logproto/logproto-indented-multiline-server.h
include/syslog-ng/logproto/logproto-record-server.h
include/syslog-ng/logproto/logproto-regexp-multiline-server.h
include/syslog-ng/logproto/logproto-server.h
include/syslog-ng/logproto/logproto-text-client.h
include/syslog-ng/logproto/logproto-text-server.h
include/syslog-ng/logproto/logproto.h
include/syslog-ng/logqueue-fifo.h
include/syslog-ng/logqueue.h
include/syslog-ng/logreader.h
include/syslog-ng/logsource.h
include/syslog-ng/logstamp.h
include/syslog-ng/logthrdestdrv.h
include/syslog-ng/logwriter.h
include/syslog-ng/mainloop-call.h
include/syslog-ng/mainloop-io-worker.h
include/syslog-ng/mainloop-worker.h
include/syslog-ng/mainloop.h
include/syslog-ng/memtrace.h
include/syslog-ng/messages.h
include/syslog-ng/misc.h
include/syslog-ng/ml-batched-timer.h
include/syslog-ng/msg-format.h
include/syslog-ng/nvtable.h
include/syslog-ng/parse-number.h
include/syslog-ng/parser/parser-expr-parser.h
include/syslog-ng/parser/parser-expr.h
include/syslog-ng/pathutils.h
include/syslog-ng/persist-state.h
include/syslog-ng/persistable-state-header.h
include/syslog-ng/persistable-state-presenter.h
include/syslog-ng/plugin-types.h
include/syslog-ng/plugin.h
include/syslog-ng/poll-events.h
include/syslog-ng/poll-fd-events.h
include/syslog-ng/pragma-parser.h
include/syslog-ng/presented-persistable-state.h
include/syslog-ng/rcptid.h
include/syslog-ng/reloc.h
include/syslog-ng/rewrite/rewrite-expr-parser.h
include/syslog-ng/rewrite/rewrite-expr.h
include/syslog-ng/rewrite/rewrite-groupset.h
include/syslog-ng/rewrite/rewrite-set-tag.h
include/syslog-ng/rewrite/rewrite-set.h
include/syslog-ng/rewrite/rewrite-subst.h
include/syslog-ng/ringbuffer.h
include/syslog-ng/run-id.h
include/syslog-ng/scratch-buffers.h
include/syslog-ng/serialize.h
include/syslog-ng/service-management.h
include/syslog-ng/stats/stats-cluster.h
include/syslog-ng/stats/stats-counter.h
include/syslog-ng/stats/stats-csv.h
include/syslog-ng/stats/stats-log.h
include/syslog-ng/stats/stats-registry.h
include/syslog-ng/stats/stats-syslog.h
include/syslog-ng/stats/stats.h
include/syslog-ng/str-format.h
include/syslog-ng/syslog-names.h
include/syslog-ng/syslog-ng.h
include/syslog-ng/tags.h
include/syslog-ng/template/templates.h
include/syslog-ng/thread-utils.h
include/syslog-ng/timeutils.h
include/syslog-ng/tls-support.h
include/syslog-ng/tlscontext.h
include/syslog-ng/transport/logtransport.h
include/syslog-ng/transport/transport-aux-data.h
include/syslog-ng/transport/transport-device.h
include/syslog-ng/transport/transport-file.h
include/syslog-ng/transport/transport-pipe.h
include/syslog-ng/transport/transport-socket.h
include/syslog-ng/transport/transport-tls.h
include/syslog-ng/type-hinting.h
include/syslog-ng/uuid.h
include/syslog-ng/value-pairs.h
include/syslog-ng/versioning.h
include/syslog-ng/vptransform.h
lib/libsyslog-ng-3.6.so.0.0.0
lib/libsyslog-ng.so
lib/libsyslog-ng-3.6.so.0
lib/syslog-ng/libafamqp.so
lib/syslog-ng/libaffile.so
lib/syslog-ng/libafmongodb.so
lib/syslog-ng/libafprog.so
lib/syslog-ng/libafsocket-notls.so
lib/syslog-ng/libafsocket.so
lib/syslog-ng/libafstomp.so
lib/syslog-ng/libafuser.so
lib/syslog-ng/libbasicfuncs.so
lib/syslog-ng/libconfgen.so
lib/syslog-ng/libcryptofuncs.so
lib/syslog-ng/libcsvparser.so
lib/syslog-ng/libdbparser.so
lib/syslog-ng/libgraphite.so
lib/syslog-ng/liblinux-kmsg-format.so
lib/syslog-ng/libpseudofile.so
lib/syslog-ng/libsyslog-ng-crypto.so
lib/syslog-ng/libsyslogformat.so
lib/syslog-ng/libsystem-source.so
lib/syslog-ng/libtest/libsyslog-ng-test.a
libdata/pkgconfig/syslog-ng.pc
libdata/pkgconfig/syslog-ng-test.pc
man/man1/loggen.1.gz
man/man1/pdbtool.1.gz
man/man1/syslog-ng-ctl.1.gz
man/man5/syslog-ng.conf.5.gz
man/man8/syslog-ng.8.gz
sbin/syslog-ng
sbin/syslog-ng-ctl
share/syslog-ng/include/scl/graphite/plugin.conf
share/syslog-ng/include/scl/graphite/README
share/syslog-ng/include/scl/nodejs/plugin.conf
share/syslog-ng/include/scl/pacct/plugin.conf
share/syslog-ng/include/scl/syslogconf/README
share/syslog-ng/include/scl/syslogconf/convert-syslogconf.awk
share/syslog-ng/include/scl/syslogconf/plugin.conf
share/syslog-ng/include/scl/system/plugin.conf
share/syslog-ng/include/scl/rewrite/cc-mask.conf
share/syslog-ng/tools/cfg-grammar.y
share/syslog-ng/tools/lex-rules.am
share/syslog-ng/tools/merge-grammar.pl
share/syslog-ng/tools/system-expand
share/syslog-ng/xsd/patterndb-1.xsd
share/syslog-ng/xsd/patterndb-2.xsd
share/syslog-ng/xsd/patterndb-3.xsd
share/syslog-ng/xsd/patterndb-4.xsd
%%PORTDOCS%%%%DOCSDIR%%/AUTHORS
%%PORTDOCS%%%%DOCSDIR%%/COPYING
%%PORTDOCS%%%%DOCSDIR%%/NEWS.md
@dir etc/patterndb.d
