# src/freenas/etc/systemd/system/nfs-idmapd.service.d/override.conf
# Disable systemd rate-limiting for NFS services. The middleware may
# restart NFS frequently during configuration changes and the default
# limit of 5 starts in 10 seconds can cause dependent services to fail.
[Unit]
StartLimitIntervalSec=0
