# src/dispatcher/etc/dispatcher
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: dispatcher
# REQUIRE: datastore

. /etc/rc.subr

name="dispatcher"
pidfile="/var/run/dispatcher.pid"
command="/usr/sbin/daemon"
start_precmd="prestart"

prestart()
{
    rc_flags="-P ${pidfile} -r /usr/local/sbin/dispatcherd ${rc_flags}"
}

load_rc_config $name
run_rc_command $*