# nas_ports/net/netatalk3/files/netatalk.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: netatalk
# REQUIRE: DAEMON %%ZEROCONF%%
# KEYWORD: shutdown
#
# AFP fileserver for Mac clients.  Add the following to /etc/rc.conf to
# enable:
#
# netatalk_enable="YES"
#

netatalk_enable=${netatalk_enable-"NO"}

. /etc/rc.subr

name=netatalk
rcvar=netatalk_enable

load_rc_config ${name}

command=/usr/local/sbin/${name}
stop_cmd="stop_netatalk"
extra_commands=reload
reload_cmd="reload_netatalk"

stop_netatalk() {
	killall netatalk
	killall afpd
	killall cnid_metad
}

reload_netatalk() {
	killall -1 netatalk
}

extra_commands="reload"
reload_cmd="netatalk_reload"

netatalk_reload()
{
        local status

        if ! status=`run_rc_command status 2>&1`; then
                echo $status
                return 1
        fi
        echo 'Reloading netatalk.'
        kill -HUP $rc_pid
}

run_rc_command "$1"
