# nas_ports/freenas/freenas-networkd/pkg-plist
/usr/local/lib/networkd/src/main.py
/usr/local/etc/rc.d/networkd
/usr/local/sbin/networkd
