# src/freenas/etc/systemd/system/fsidd.service.d/override.conf
[Unit]
# stop fsidd.service when nfs-server stops
PartOf=nfs-server.service
