# nas_ports/freenas/pcbsd-warden/pkg-plist
usr/local/bin/warden
usr/local/share/warden/pluginjail-packages
usr/local/share/warden/pcbsd-utils-packages
usr/local/etc/warden.conf
usr/local/share/warden/linux-installs/debian-6-squeeze
usr/local/share/warden/linux-installs/gentoo-stage3-i486
usr/local/share/pcbsd/scripts/functions.sh
usr/local/share/warden/scripts/backend/checkstatus.sh
usr/local/share/warden/scripts/backend/checkupdates.sh
usr/local/share/warden/scripts/backend/chrootjail.sh
usr/local/share/warden/scripts/backend/createjail.sh
usr/local/share/warden/scripts/backend/createtemplate.sh
usr/local/share/warden/scripts/backend/cronsnap.sh
usr/local/share/warden/scripts/backend/deletejail.sh
usr/local/share/warden/scripts/backend/details.sh
usr/local/share/warden/scripts/backend/exportjail.sh
usr/local/share/warden/scripts/backend/functions.sh
usr/local/share/warden/scripts/backend/importjail.sh
usr/local/share/warden/scripts/backend/jailid.sh
usr/local/share/warden/scripts/backend/listjails.sh
usr/local/share/warden/scripts/backend/listpbis.sh
usr/local/share/warden/scripts/backend/listpkgs.sh
usr/local/share/warden/scripts/backend/setupusers.sh
usr/local/share/warden/scripts/backend/startjail.sh
usr/local/share/warden/scripts/backend/stopjail.sh
usr/local/share/warden/scripts/backend/zfsclonesnap.sh
usr/local/share/warden/scripts/backend/zfscronsnap.sh
usr/local/share/warden/scripts/backend/zfslistclone.sh
usr/local/share/warden/scripts/backend/zfslistsnap.sh
usr/local/share/warden/scripts/backend/zfsmksnap.sh
usr/local/share/warden/scripts/backend/zfsrevertsnap.sh
usr/local/share/warden/scripts/backend/zfsrmclone.sh
usr/local/share/warden/scripts/backend/zfsrmsnap.sh
usr/local/share/warden/scripts/hooks/dhclient-exit-hooks
usr/local/share/warden/scripts/rc.d/shutdown.sh
usr/local/share/warden/scripts/rc.d/startup.sh
usr/local/share/warden/scripts/rc.d/wardenrc
usr/local/etc/wardenrc
usr/local/share/warden/agent.png
usr/local/share/warden/warden.png
usr/local/share/pcbsd/conf/iso-mirror
usr/local/share/pcbsd/conf/pbi-mirror
usr/local/share/pcbsd/conf/pcbsd-mirrors
usr/local/share/pcbsd/conf/pkg-mirror
usr/local/share/pcbsd/conf/update-mirror
