# src/freenas/etc/systemd/system/nvmf-connect@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/sh -c "/usr/local/bin/nvmf-wrapper.sh `/bin/echo -e '${CONNECT_ARGS}'`"
