# nas_ports/freenas/freenas-etcd/pkg-plist
/usr/local/lib/etcd/src/renderers.py
/usr/local/lib/etcd/src/main.py
/usr/local/lib/etcd/plugins/group.mako
/usr/local/lib/etcd/plugins/loader.conf.mako
/usr/local/lib/etcd/plugins/local/nginx/nginx.conf.mako
/usr/local/lib/etcd/plugins/local/collectd.conf.mako
/usr/local/lib/etcd/plugins/master.passwd.mako
/usr/local/lib/etcd/plugins/passwd.mako
/usr/local/lib/etcd/plugins/rc.conf.mako
/usr/local/lib/etcd/plugins/resolv.conf.mako
/usr/local/lib/etcd/plugins/pwd_database.py
/usr/local/lib/etcd/plugins/pubkeys.py
/usr/local/lib/etcd/plugins/ssh/sshd_config.mako
/usr/local/lib/etcd/plugins/ssh/sshd_keys.py
/usr/local/lib/etcd/plugins/exports.mako
/usr/local/etc/rc.d/etcd
/usr/local/sbin/etcd
