# nas_ports/www/minio/files/minio.in
#!/bin/sh

# PROVIDE: minio
# KEYWORD: shutdown

minio_enable="${minio_enable-NO}"
minio_config="${minio_config-%%PREFIX%%/etc/minio}"
minio_env=${minio_env:=""}
minio_disks="${minio_disks}"
minio_address="${minio_address-:9000}"

. /etc/rc.subr

name=minio
procname="%%PREFIX%%/bin/minio"
rcvar=minio_enable

load_rc_config ${name}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-c -p ${pidfile} /usr/local/bin/${name} -C \"${minio_config}\" server --address=\"${minio_address}\" --quiet ${minio_disks}"

start_precmd="minio_startprecmd"

minio_startprecmd()
{
	if [ ! -e "${pidfile}" ]; then
		install "/dev/null" "${pidfile}"
	fi

	if [ ! -d "${minio_config}" ]; then
		mkdir -p "${minio_config}"
	fi
}

run_rc_command "$1"
