# src/freenas/etc/systemd/system/kdump-tools-dump.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/local/bin/execute_dump.sh
