# nas_ports/sysutils/stanchion/files/stanchion.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: stanchion
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name=stanchion
command=/usr/local/lib/stanchion/erts-5.10.3/bin/beam.smp
rcvar=stanchion_enable
start_cmd="/usr/local/sbin/${name} start"
stop_cmd="/usr/local/sbin/${name} stop"
pidfile="/var/run/${name}/${name}.pid"

load_rc_config $name
run_rc_command "$1"
