# nas_ports/databases/riak2/files/riak.in
#!/bin/sh                                                                                                               

# $FreeBSD$
#
# PROVIDE: riak
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# riak_enable (bool):   Set to NO by default.
#                       Set it to YES to enable riak

. /etc/rc.subr

name=riak
rcvar=riak_enable

load_rc_config $name

command=%%PREFIX%%/lib/riak/erts-5.10.4/bin/beam.smp
pidfile="/var/run/${name}/${name}.pid"
start_cmd="%%PREFIX%%/sbin/${name} start"
stop_cmd="%%PREFIX%%/sbin/${name} stop"

run_rc_command "$1"
