# src/freenas/etc/systemd/system/serial-getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/cli --menu
StandardInput=tty
StandardOutput=tty
