# src/freenas/etc/systemd/system/serial-getty@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/etc/netcli
StandardInput=tty
StandardOutput=tty
