# nas_ports/freenas/freenas-dispatcher/pkg-plist
/usr/local/lib/python2.7/site-packages/dispatcher_client-10.0-py2.7.egg-info
/usr/local/lib/python2.7/site-packages/dispatcher/rpc.pyc
/usr/local/lib/python2.7/site-packages/dispatcher/rpc.py
/usr/local/lib/python2.7/site-packages/dispatcher/client.py
/usr/local/lib/python2.7/site-packages/dispatcher/__init__.pyc
/usr/local/lib/python2.7/site-packages/dispatcher/__init__.py
/usr/local/lib/python2.7/site-packages/dispatcher/client.pyc
/usr/local/lib/python2.7/site-packages/dispatcher/validator.py
/usr/local/lib/python2.7/site-packages/dispatcher/validator.pyc
/usr/local/lib/python2.7/site-packages/dispatcher/jsonenc.py
/usr/local/lib/python2.7/site-packages/dispatcher/jsonenc.pyc
/usr/local/lib/python2.7/site-packages/dispatcher/shell.py
/usr/local/lib/python2.7/site-packages/dispatcher/shell.pyc
/usr/local/lib/dispatcher/src/main.py
/usr/local/lib/dispatcher/src/api/__init__.py
/usr/local/lib/dispatcher/src/api/handler.py
/usr/local/lib/dispatcher/src/task.py
/usr/local/lib/dispatcher/src/event.py
/usr/local/lib/dispatcher/src/cache.py
/usr/local/lib/dispatcher/src/query.py
/usr/local/lib/dispatcher/src/lib/freebsd.py
/usr/local/lib/dispatcher/src/lib/geom.py
/usr/local/lib/dispatcher/src/lib/__init__.py
/usr/local/lib/dispatcher/src/lib/system.py
/usr/local/lib/dispatcher/src/balancer.py
/usr/local/lib/dispatcher/src/auth.py
/usr/local/lib/dispatcher/src/services.py
/usr/local/lib/dispatcher/src/resources.py
/usr/local/lib/dispatcher/src/utils.py
/usr/local/lib/dispatcher/src/frontend/templates/nav.html
/usr/local/lib/dispatcher/src/frontend/templates/tasks.html
/usr/local/lib/dispatcher/src/frontend/templates/events.html
/usr/local/lib/dispatcher/src/frontend/templates/rpc.html
/usr/local/lib/dispatcher/src/frontend/templates/term.html
/usr/local/lib/dispatcher/src/frontend/templates/stats.html
/usr/local/lib/dispatcher/src/frontend/templates/apidoc/index.html
/usr/local/lib/dispatcher/src/frontend/templates/apidoc/nav.html
/usr/local/lib/dispatcher/src/frontend/templates/apidoc/rpc.html
/usr/local/lib/dispatcher/src/frontend/templates/apidoc/events.html
/usr/local/lib/dispatcher/src/frontend/templates/apidoc/tasks.html
/usr/local/lib/dispatcher/src/frontend/templates/apidoc/schemas.html
/usr/local/lib/dispatcher/src/frontend/static/middleware.js
/usr/local/lib/dispatcher/src/frontend/static/emitter.js
/usr/local/lib/dispatcher/src/frontend/static/term.js
/usr/local/lib/dispatcher/src/frontend/static/jquery.jsonview.js
/usr/local/lib/dispatcher/src/frontend/static/jquery.jsonview.css
/usr/local/lib/dispatcher/src/frontend/static/main.css
/usr/local/lib/dispatcher/src/frontend/frontend.py
/usr/local/lib/dispatcher/src/frontend/__init__.py
/usr/local/lib/dispatcher/plugins/SystemInfoPlugin.py
/usr/local/lib/dispatcher/plugins/ZfsPlugin.py
/usr/local/lib/dispatcher/plugins/DevdPlugin.py
/usr/local/lib/dispatcher/plugins/ServiceManagePlugin.py
/usr/local/lib/dispatcher/plugins/ConfigPlugin.py
/usr/local/lib/dispatcher/plugins/AuditPipePlugin.py
/usr/local/lib/dispatcher/plugins/UserPlugin.py
/usr/local/lib/dispatcher/plugins/FilesystemPlugin.py
/usr/local/lib/dispatcher/plugins/VolumePlugin.py
/usr/local/lib/dispatcher/plugins/DiskPlugin.py
/usr/local/lib/dispatcher/plugins/NetworkPlugin.py
/usr/local/lib/dispatcher/plugins/EntitySubscriberPlugin.py
/usr/local/lib/dispatcher/plugins/SessionPlugin.py
/usr/local/lib/dispatcher/plugins/ShareNFSPlugin.py
/usr/local/lib/dispatcher/plugins/ShareAFPPlugin.py
/usr/local/lib/dispatcher/plugins/SharingPlugin.py
/usr/local/lib/dispatcher/plugins/SystemDatasetPlugin.py
/usr/local/lib/dispatcher/plugins/UpdatePlugin.py
/usr/local/libexec/dispatcher/notify-pam
/usr/local/etc/rc.d/dispatcher
/usr/local/sbin/dispatcherd
/usr/local/sbin/dispatcherctl
/usr/local/etc/middleware.conf

