# src/pcbsd/pbi-manager/port-files/pkg-plist
sbin/pbi-crashhandler
sbin/pbi_add
sbin/pbi_addrepo
sbin/pbi_autobuild
sbin/pbi_browser
sbin/pbi_create
sbin/pbi_delete
sbin/pbi_deleterepo
sbin/pbi_icon
sbin/pbi_indextool
sbin/pbi_info
sbin/pbi_listrepo
sbin/pbi_makepatch
sbin/pbi_makeport
sbin/pbi_makerepo
sbin/pbi_metatool
sbin/pbi_update
sbin/pbi_update_hashdir
sbin/pbi_patch
sbin/pbid
etc/pbi.conf
etc/rc.d/pbid
%%DATADIR%%/module-examples/gimp/xdg-menu/gimp.desktop
%%DATADIR%%/module-examples/gimp/xdg-mime/gimp-xdg.png
%%DATADIR%%/module-examples/gimp/xdg-mime/gimp-xdg.xml
%%DATADIR%%/module-examples/gimp/LICENSE
%%DATADIR%%/module-examples/gimp/external-links
%%DATADIR%%/module-examples/gimp/pbi.conf
%%DATADIR%%/module-examples/gimp/scripts/pre-remove.sh
%%DATADIR%%/module-examples/gimp/scripts/post-install.sh
%%DATADIR%%/module-examples/gimp/scripts/pre-portmake.sh
%%DATADIR%%/module-examples/gimp/scripts/post-portmake.sh
%%DATADIR%%/module-examples/gimp/scripts/pre-install.sh
%%DATADIR%%/module-examples/gimp/resources/gui_sidebanner.png
%%DATADIR%%/module-examples/gimp/resources/gimp.png
%%DATADIR%%/module-examples/gimp/resources/gui_banner.png
%%DATADIR%%/module-examples/gimp/xdg-desktop/gimp.desktop
%%DATADIR%%/module-examples/convertoldmod.sh
%%DATADIR%%/module-examples/firefox/LICENSE
%%DATADIR%%/module-examples/firefox/resources/gui_banner.png
%%DATADIR%%/module-examples/firefox/resources/firefox.png
%%DATADIR%%/module-examples/firefox/resources/gui_sidebanner.png
%%DATADIR%%/module-examples/firefox/xdg-menu/firefox3.desktop
%%DATADIR%%/module-examples/firefox/scripts/pre-remove.sh
%%DATADIR%%/module-examples/firefox/scripts/post-install.sh
%%DATADIR%%/module-examples/firefox/scripts/pre-install.sh
%%DATADIR%%/module-examples/firefox/scripts/post-portmake.sh
%%DATADIR%%/module-examples/firefox/xdg-desktop/firefox.desktop
%%DATADIR%%/module-examples/firefox/external-links
%%DATADIR%%/module-examples/firefox/pbi.conf
%%DATADIR%%/icons/default.png
%%DATADIR%%/icons/patch.png
@dirrm %%DATADIR%%/module-examples/firefox/resources
@dirrm %%DATADIR%%/module-examples/firefox/scripts
@dirrm %%DATADIR%%/module-examples/firefox/xdg-desktop
@dirrm %%DATADIR%%/module-examples/firefox/xdg-menu
@dirrm %%DATADIR%%/module-examples/firefox/xdg-mime
@dirrm %%DATADIR%%/module-examples/firefox
@dirrm %%DATADIR%%/module-examples/gimp/resources
@dirrm %%DATADIR%%/module-examples/gimp/scripts
@dirrm %%DATADIR%%/module-examples/gimp/xdg-desktop
@dirrm %%DATADIR%%/module-examples/gimp/xdg-menu
@dirrm %%DATADIR%%/module-examples/gimp/xdg-mime
@dirrm %%DATADIR%%/module-examples/gimp
@dirrm %%DATADIR%%/module-examples
@dirrm %%DATADIR%%/icons
@dirrm %%DATADIR%%
@unexec rmdir /usr/pbi/.hashdir >/dev/null 2>/dev/null
@unexec rmdir /usr/pbi >/dev/null 2>/dev/null
