# src/freenas/etc/systemd/system/syslog-ng.service.d/override.conf
[Unit]
# Exclude syslog-ng from the restart limiter.
StartLimitIntervalSec=0
