# nas_ports/freenas/py-middlewared/files/middlewared.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: middlewared
# BEFORE: fsck earlykld

. /etc/rc.subr

middlewared_start() {
	/sbin/ifconfig lo0 127.0.0.1 up
	${command}
	# FIXME: add midctl command to wait for server to be ready
	# for now, the old and brave sleep
	sleep 5
}

name="middlewared"
pidfile="/var/run/middlewared.pid"
command="/usr/local/bin/${name}"
start_cmd='middlewared_start'

load_rc_config $name
run_rc_command $*
