# nas_ports/freenas/py-middlewared/files/middlewared.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: middlewared
# BEFORE: fsck earlykld
# KEYWORD: shutdown

. /etc/rc.subr

middlewared_debug=${middlewared_debug:-"NO"}
pidfile=/var/run/middlewared.pid
command="/usr/sbin/daemon"

middlewared_start() {
	mount -uw /  # FIXME: why root is not rw by default?
	/sbin/ifconfig lo0 127.0.0.1 up
	if checkyesno middlewared_debug; then
		/usr/local/bin/tmux new-session -s middlewared -d
		/usr/local/bin/tmux send -t middlewared "env LC_ALL=en_US.UTF-8 ${command} -f" ENTER
	else
		env PATH=$PATH:/usr/local/sbin:/usr/local/bin LC_ALL=en_US.UTF-8 ${command} -f -P ${pidfile} -r /usr/local/bin/middlewared -f --log-handler=file
	fi
	LD_LIBRARY_PATH=/usr/local/lib /usr/local/bin/midclt -t 30 waitready
}

name="middlewared"
start_cmd='middlewared_start'

load_rc_config $name
run_rc_command $*
