# src/freenas/etc/systemd/system/dbus.socket.d/override.conf
[Unit]
# It implicitly depends on sysinit.target which depends on local-fs.target
# But we start middleware before local-fs, but middleware needs D-Bus
DefaultDependencies=no
