# src/freenas/etc/systemd/system/nfs-server.service.d/override.conf
[Unit]
ConditionDirectoryNotEmpty=!/etc/exports.d
