# src/freenas/etc/systemd/system/nfs-server.service.d/override.conf
[Unit]
ConditionDirectoryNotEmpty=!/etc/exports.d
After=nslcd.service

[Service]
ExecStartPre=/usr/local/bin/clear_nfs_rmtab.sh
