# nas_ports/net-mgmt/collectd5/files/collectd.in
#!/bin/sh

# PROVIDE: collectd
# REQUIRE: LOGIN
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# collectd_enable (bool):       Set to NO by default.
#               Set it to YES to enable collectd.
# collectd_flags (string):      Set to "" by default.
#               Additional flags to collectd, see manual page.
#

. /etc/rc.subr

name="collectd"
rcvar=collectd_enable

command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${collectd_enable="NO"}

run_rc_command "$1"
