# nas_ports/net-mgmt/collectd5/files/collectdmon.in
#!/bin/sh

# PROVIDE: collectdmon
# REQUIRE: LOGIN
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# collectdmon_enable (bool):   Set to NO by default.
#               Set it to YES to enable collectd.
# collectdmon_flags (string):  Set to "-c %%PREFIX%%/sbin/collectd" by default.
#               Additional flags to collectdmon, see manual page.
#

. /etc/rc.subr

name="collectdmon"
rcvar=collectdmon_enable

command=%%PREFIX%%/sbin/${name}

load_rc_config $name

: ${collectdmon_enable="NO"}
: ${collectdmon_flags="-c %%PREFIX%%/sbin/collectd"}

run_rc_command "$1"
