# nas_ports/freenas/freenas-etcfs/pkg-plist
/usr/local/lib/etcfs/src/renderers.py
/usr/local/lib/etcfs/src/main.py
/usr/local/lib/etcfs/src/fs.py
/usr/local/lib/etcfs/plugins/group.mako
/usr/local/lib/etcfs/plugins/loader.conf.mako
/usr/local/lib/etcfs/plugins/local/nginx.conf.mako
/usr/local/lib/etcfs/plugins/master.passwd.mako
/usr/local/lib/etcfs/plugins/passwd.mako
/usr/local/lib/etcfs/plugins/rc.conf.mako
/usr/local/lib/etcfs/plugins/resolv.conf.mako
/usr/local/etc/rc.d/etcfs
/usr/local/sbin/mount_etcfs
