# nas_ports/freenas/freenas-etcfs/pkg-plist
/usr/local/lib/etcfs/renderers.py
/usr/local/lib/etcfs/main.py
/usr/local/lib/etcfs/fs.py
/usr/local/etc/rc.d/etcfs
/usr/local/sbin/mount_etcfs
