# nas_ports/freenas/freenas-files/pkg-plist
usr/local/bin/extract-tarball
usr/local/bin/winacl
usr/local/bin/autotune
usr/local/bin/afpusers
usr/local/bin/arc_summary.py
usr/local/bin/arcstat.py
usr/local/bin/crashuploader
usr/local/bin/custom-upssched-cmd
usr/local/bin/fixsid.py
usr/local/bin/freenas-debug
usr/local/bin/ixdiagnose
usr/local/bin/mfistatus.py
usr/local/bin/perftests-nas
usr/local/bin/sas2ircu
usr/local/bin/zilstat
usr/local/bin/adtool
usr/local/bin/ldaptool
usr/local/bin/chelsio_adapter_config
usr/local/bin/t5seeprom
usr/local/bin/sockscopy
usr/local/bin/freenas-snmp/__init__.py
usr/local/bin/freenas-snmp/snmpresponse.py
usr/local/bin/freenas-snmp/zfs-snmp
usr/local/etc/avahi/avahi-daemon.conf
usr/local/etc/default/grub
usr/local/etc/nginx/fastcgi_params
usr/local/etc/nginx/nginx.conf
usr/local/etc/rc.d/django
usr/local/etc/rc.d/ldpbiconfig
usr/local/etc/rc.d/pbid
usr/local/etc/samba.gdb
usr/local/etc/sudoers
usr/local/etc/syslog-ng.conf.freenas
usr/local/libexec/freenas-debug/active_directory/active_directory.sh
usr/local/libexec/freenas-debug/cache/cache.sh
usr/local/libexec/freenas-debug/domain_controller/domain_controller.sh
usr/local/libexec/freenas-debug/dtrace/disklatencycmd.d
usr/local/libexec/freenas-debug/dtrace/dtrace.sh
usr/local/libexec/freenas-debug/geom/geom.sh
usr/local/libexec/freenas-debug/hardware/hardware.sh
usr/local/libexec/freenas-debug/include.sh
usr/local/libexec/freenas-debug/ipmi/ipmi.sh
usr/local/libexec/freenas-debug/iscsi/iscsi.sh
usr/local/libexec/freenas-debug/ldap/ldap.sh
usr/local/libexec/freenas-debug/loader/loader.sh
usr/local/libexec/freenas-debug/network/network.sh
usr/local/libexec/freenas-debug/nfs/nfs.sh
usr/local/libexec/freenas-debug/smart/smart.sh
usr/local/libexec/freenas-debug/ssl/ssl.sh
usr/local/libexec/freenas-debug/sysctl/sysctl.sh
usr/local/libexec/freenas-debug/system/system.sh
usr/local/libexec/freenas-debug/zfs/zfs.sh
usr/local/libexec/freenas-debug/verify/verify.sh
usr/local/libexec/nas/generate_afpd_conf.py
usr/local/libexec/nas/generate_ctl_conf.py
usr/local/libexec/nas/generate_krb5_conf.py
usr/local/libexec/nas/generate_ldap_conf.py
usr/local/libexec/nas/generate_nsswitch_conf.py
usr/local/libexec/nas/generate_resolv_conf.py
usr/local/libexec/nas/generate_smb4_conf.py
usr/local/libexec/nas/generate_ssl_certificates.py
usr/local/libexec/nas/generate_sssd_conf.py
usr/local/libexec/nas/generate_webdav_conf.py
usr/local/libexec/nas/register_mdns.py
usr/local/libexec/samba/dfree
usr/local/libexec/samba/samba-backtrace
usr/local/sbin/sas2flash
usr/local/sbin/save_rrds.sh
usr/local/sbin/scrub
usr/local/share/firmware/t5fw-1.12.25.0.bin
usr/local/share/perftests-nas/tests/iozone
usr/local/share/vpds/t580_cr_qsa_variable_2133_vpd.bin
usr/local/share/vpds/t580_cr_spider_variable_2133_vpd.bin
usr/local/share/vpds/t580_cr_variable_2133_vpd.bin
usr/local/share/vpds/t580_lp_cr_qsa_variable_2133_vpd.bin
usr/local/share/vpds/t580_lp_cr_spider_variable_2133_vpd.bin
usr/local/share/vpds/t580_lp_cr_variable_2133_vpd.bin
usr/local/share/vpds/t580_lp_so_qsa_variable_vpd.bin
usr/local/share/vpds/t580_lp_so_spider_variable_2133_vpd.bin
usr/local/share/vpds/t580_lp_so_variable_vpd.bin
usr/local/share/snmp/mibs/FREENAS-MIB.txt
usr/local/www/data/django.fcgi
etc/ix.rc.d/ix-activedirectory
etc/ix.rc.d/ix-afpd
etc/ix.rc.d/ix-aliases
etc/ix.rc.d/ix-ataidle
etc/ix.rc.d/ix-snmpd
etc/ix.rc.d/ix-cache
etc/ix.rc.d/ix-collectd
etc/ix.rc.d/ix-crontab
etc/ix.rc.d/ix-ctld
etc/ix.rc.d/ix-device_hints
etc/ix.rc.d/ix-fstab
etc/ix.rc.d/ix-hostname
etc/ix.rc.d/ix-inadyn
etc/ix.rc.d/ix-inetd
etc/ix.rc.d/ix-jail
etc/ix.rc.d/ix-kerberos
etc/ix.rc.d/ix-kinit
etc/ix.rc.d/ix-ldap
etc/ix.rc.d/ix-loader
etc/ix.rc.d/ix-localtime
etc/ix.rc.d/ix-motd
etc/ix.rc.d/ix-multipath
etc/ix.rc.d/ix-nfsd
etc/ix.rc.d/ix-nginx
etc/ix.rc.d/ix-nis
etc/ix.rc.d/ix-nsswitch
etc/ix.rc.d/ix-nt4
etc/ix.rc.d/ix-ntpd
etc/ix.rc.d/ix-pam
etc/ix.rc.d/ix-passwd
etc/ix.rc.d/ix-plugins
etc/ix.rc.d/ix-postinit
etc/ix.rc.d/ix-preinit
etc/ix.rc.d/ix-proftpd
etc/ix.rc.d/ix-resolv
etc/ix.rc.d/ix-rsyncd
etc/ix.rc.d/ix-pre-samba
etc/ix.rc.d/ix-post-samba
etc/ix.rc.d/ix-savehostid
etc/ix.rc.d/ix-shutdown
etc/ix.rc.d/ix-smartd
etc/ix.rc.d/ix-smbpasswd
etc/ix.rc.d/ix-sshd
etc/ix.rc.d/ix-ssl
etc/ix.rc.d/ix-sssd
etc/ix.rc.d/ix-sudoers
etc/ix.rc.d/ix-swap
etc/ix.rc.d/ix-syncdisks
etc/ix.rc.d/ix-sysctl
etc/ix.rc.d/ix-syslogd
etc/ix.rc.d/ix-system
etc/ix.rc.d/ix-ttys
etc/ix.rc.d/ix-update
etc/ix.rc.d/ix-ups
etc/ix.rc.d/ix-warden
etc/ix.rc.d/ix-zfs
etc/ix.rc.d/ix_diagnose
etc/ix.rc.d/ix_pf_early
etc/ix.rc.d/ix_pf_late
etc/ix.rc.d/ix_register
etc/ix.rc.d/ix_sshd_save_keys
etc/ix.rc.d/ix_textdump
etc/ix.rc.d/ix-apache
etc/directoryservice/ActiveDirectory/ctl
etc/directoryservice/DomainController/ctl
etc/directoryservice/LDAP/ctl
etc/directoryservice/NIS/ctl
etc/directoryservice/NT4/ctl
etc/directoryservice/rc.ActiveDirectory
etc/directoryservice/rc.DomainController
etc/directoryservice/rc.LDAP
etc/directoryservice/rc.NIS
etc/directoryservice/rc.NT4
etc/ix/templates/kerberos/krb5.conf
etc/ix/templates/pam.d/README
etc/ix/templates/pam.d/atrun
etc/ix/templates/pam.d/cron
etc/ix/templates/pam.d/ftp
etc/ix/templates/pam.d/ftpd
etc/ix/templates/pam.d/imap
etc/ix/templates/pam.d/kde
etc/ix/templates/pam.d/login
etc/ix/templates/pam.d/netatalk
etc/ix/templates/pam.d/other
etc/ix/templates/pam.d/passwd
etc/ix/templates/pam.d/pop3
etc/ix/templates/pam.d/rsh
etc/ix/templates/pam.d/samba
etc/ix/templates/pam.d/sshd
etc/ix/templates/pam.d/su
etc/ix/templates/pam.d/sudo
etc/ix/templates/pam.d/system
etc/ix/templates/pam.d/telnetd
etc/ix/templates/pam.d/xdm
etc/iso_3166_2_countries.csv
etc/netcli
@shell etc/netcli.sh
etc/find_alias_for_smtplib.py
