# nas_ports/freenas/freenas-files/pkg-plist
boot/loader.conf
usr/local/bin/extract-tarball
usr/local/bin/winacl
usr/local/bin/autotune
usr/local/bin/afpusers
usr/local/bin/arc_summary.py
usr/local/bin/arcstat.py
usr/local/bin/crashuploader
usr/local/bin/custom-upssched-cmd
usr/local/bin/ixdiagnose
usr/local/bin/mfistatus.py
usr/local/bin/perftests-nas
usr/local/bin/zilstat
usr/local/bin/adtool
usr/local/etc/avahi/avahi-daemon.conf
usr/local/etc/appcafe.conf
usr/local/etc/samba.gdb
usr/local/etc/sudoers
usr/local/etc/warden.conf
usr/local/libexec/samba/samba-backtrace
usr/local/sbin/sas2flash
usr/local/sbin/save_rrds.sh
usr/local/sbin/scrub
usr/local/share/perftests-nas/tests/iozone
etc/local
etc/csh.cshrc
etc/devfs.rules
etc/exports
etc/freenas.conf
etc/gettytab
etc/ipfw.conf.block
etc/newsyslog.conf
etc/pf.conf.block
etc/pw.conf
etc/rc.subr
etc/rc.conf
etc/sysctl.conf
etc/ttys
etc/mail/mailer.conf
etc/rc.d/cleanvar
etc/rc.d/ix-bootstrap
etc/rc.d/ix-etc
etc/rc.d/ix-network
etc/launchd.d/datastore.json
etc/launchd.d/dispatcher.json
.bashrc
.cshrc
.profile
.shrc
@exec mkdir -p /data/zfs
