# debian/debian/ix-cgroups.service
[Unit]
Description=Enable required cgroups for TrueNAS
DefaultDependencies=no
# We want to set it up before starting middleware because of NAS-120085
Before=middlewared.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/setup_cgroups
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
