# debian/debian/ix-ssh-keys.service
[Unit]
Description=Save SSH keys
Wants=ssh.service

[Service]
ExecStart=midclt call ssh.save_keys

[Install]
WantedBy=multi-user.target
