# debian/debian/ix-vendor.service
[Unit]
# vendor_service.py must run after network routes are set up.
After=ix-netif.service


[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=30
ExecStart=/usr/bin/start_vendor_service
StandardOutput=null
