# debian/debian/ix-vendor.service
[Unit]
# Vendors can do a bunch of different things
# so we want to run it as late as possible
# in the boot process.
After=ix-postinit.service

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=30
RestartSec=10
ExecStart=-/usr/bin/start_vendor_service
StandardOutput=null

[Install]
WantedBy=truenas.target
