# debian/debian/postinst
#!/bin/sh

# It unconfigures interfaces we've just configured even if it is instructed to do not touch anything
systemctl disable systemd-networkd
