# debian/debian/postinst
#!/bin/sh -ex

mkdir -p /var/trash
for file in \
    /etc/cron.d/certbot \
    /etc/cron.d/e2scrub_all \
    /etc/cron.d/mdadm \
    /etc/cron.d/sysstat \
    /etc/cron.daily/apt-compat \
    /etc/cron.daily/apache2 \
    /etc/cron.daily/calendar \
    /etc/cron.daily/exim4-base \
    /etc/cron.daily/dpkg \
    /etc/cron.daily/mdadm \
    /etc/cron.daily/man-db \
    /etc/cron.daily/sysstat \
    /etc/cron.weekly/man-db \
    /etc/grub.d/10_linux \
    /etc/grub.d/30_os-prober \
    /etc/logrotate.d/syslog-ng
do
    dpkg-divert --add --package truenas --rename --divert "/var/trash/$(echo "$file" | sed "s/\//_/g")" "$file"
done

# openssh-server explicitly sets itself to enabled on post install, so preset is not respected
# we need to manually disable it

systemctl disable ssh.service

# Services whose unit files are created with systemd-sysv-generator do not respect preset file so for now
# we disable them explicitly

systemctl disable scst
systemctl disable inadyn
systemctl disable tftpd-hpa
systemctl disable proftpd
systemctl disable nmbd
systemctl disable smbd
systemctl disable winbind
systemctl disable wsdd
systemctl disable walinuxagent
systemctl disable openipmi

# nvidia-persistenced is not respecting vendor preset file so we disable it explicitly
systemctl disable nvidia-persistenced

# Update alternatives
update-alternatives --install "/usr/sbin/sendmail" sendmail "/etc/find_alias_for_smtplib.sh" "10"

# Add nut to dialout group - NAS-110578
usermod -a -G dialout nut
# Usbhid-ups driver does not start and needs nut as the group for it's executable to correctly execute NAS-113642
chown root:nut /lib/nut/usbhid-ups

# We would like to disable docker-compose to prevent people from foot-shooting
chmod -x /usr/bin/docker-compose

# Copy to /conf/base
for bit in /etc/aliases /etc/group /etc/passwd /etc/syslog-ng/syslog-ng.conf; do
    mkdir -p "$(dirname "/conf/base/$bit")"
    cp -aR "$bit" "/conf/base/$bit"
done

#DEBHELPER#

# zfs-zed is not respecting vendor preset file so we enable it explicitly
systemctl daemon-reload
systemctl enable zfs-zed

# We need to mask certain services so that they don't start automatically
systemctl mask libvirtd.socket libvirtd-ro.socket libvirtd-admin.socket libvirtd-tls.socket libvirtd-tcp.socket
systemctl mask exim4-base.service exim4.service exim4-base.timer
systemctl mask uuidd.service uuidd.socket
systemctl mask ndctl-monitor.service

# We don't use LVM and this service can add significant boot delays
# on large disk systems
systemctl mask lvm2-monitor.service

systemctl set-default truenas.target

sed -i.bak 's/CHARMAP="ISO-8859-15"/CHARMAP="UTF-8"/' /etc/default/console-setup
