# debian/debian/truenas.target
[Unit]
Description=TrueNAS Target
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
