# nas_ports/sysutils/nut/files/nut_upslog.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nut_upslog
# REQUIRE: NETWORKING nut
# BEFORE: LOGIN
# KEYWORD: shutdown

# Define these nut_upslog* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/nut_upslog
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
nut_upslog_enable=${nut_upslog_enable-"NO"}
nut_upslog_prefix=${nut_upslog_prefix-"%%PREFIX%%"}
nut_upslog_logfile=${nut_upslog_logfile-"/var/log/ups.log"}
nut_upslog_interval=${nut_upslog_interval-"300"}
nut_upslog_ups=${nut_upslog_ups-"myups@localhost"}

. /etc/rc.subr

name="nut_upslog"
rcvar=nut_upslog_enable

load_rc_config $name

required_dirs="%%STATEDIR%%"
command="${nut_upslog_prefix}/bin/upslog"
pidfile="%%STATEDIR%%/upslog.pid"
nut_upslog_flags=${nut_upslog_flags-"-s ${nut_upslog_ups} -l ${nut_upslog_logfile} -i ${nut_upslog_interval}"}

run_rc_command "$1"
