# nas_ports/sysutils/nut/files/nut_upsmon.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: nut_upsmon
# REQUIRE: NETWORKING nut
# BEFORE: LOGIN
# KEYWORD: shutdown

# Define these nut_upsmon* variables in one of these files:
#       /etc/rc.conf
#       /etc/rc.conf.local
#       /etc/rc.conf.d/nut_upsmon
#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
#
nut_upsmon_enable=${nut_upsmon_enable-"NO"}
nut_upsmon_prefix=${nut_upsmon_prefix-"%%PREFIX%%"}

. /etc/rc.subr

name="nut_upsmon"
rcvar=nut_upsmon_enable

load_rc_config $name

required_dirs="%%STATEDIR%%"
required_files="${nut_upsmon_prefix}/etc/nut/upsmon.conf"
command="${nut_upsmon_prefix}/sbin/upsmon"
pidfile="%%STATEDIR%%/upsmon.pid"
nut_upsmon_flags=${nut_upsmon_flags-"localhost"}

extra_commands=reload
reload()
{
	${command} -c reload
}

run_rc_command "$1"
