# nas_ports/emulators/virtualbox-ose-kmod/files/vboxnet.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE:	vboxnet
# BEFORE:	NETWORKING
# REQUIRE:	FILESYSTEMS
# KEYWORD:	nojail

#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# vboxnet_enable (bool):   Set to NO by default.
#               Set it to YES to load network related kernel modules on startup

. /etc/rc.subr

name="vboxnet"
rcvar=vboxnet_enable
start_cmd="vboxnet_start"
stop_cmd="vboxnet_stop"

vboxnet_start()
{
	if ! kldstat -q -m ng_vboxnetflt;
	then
		if ! kldload vboxnetflt > /dev/null 2>&1;
		then
			warn "Can't load vboxnetflt module."
			return 1
		fi
	fi

	if ! kldstat -q -m vboxnetadp;
	then
		if ! kldload vboxnetadp > /dev/null 2>&1;
		then
			warn "Can't load vboxnetadp module."
			return 1
		fi
	fi

	# initialize configured host-only interfaces
	# Only do this if running in a jail
	# On FreeNAS the kernel modules need to be loaded
	# from the host environment, but the host only
	# interfaces are created in the vbox jail.
	jail=$(sysctl -n security.jail.jailed)
	if [ ${jail} -ne 0 ]; then
		LD_LIBRARY_PATH=%%PREFIX%%/lib %%PREFIX%%/bin/VBoxManage list hostonlyifs >/dev/null
	fi
}

vboxnet_stop()
{
	if kldstat -q -m vboxnetadp;
	then
		if ! kldunload vboxnetadp > /dev/null 2>&1;
		then
			warn "Can't unload vboxnetadp module."
			return 1
		fi
	fi

	if kldstat -q -m ng_vboxnetflt;
	then
		if ! kldunload vboxnetflt > /dev/null 2>&1;
		then
			warn "Can't unload vboxnetflt module."
			return 1
		fi
	fi
}

load_rc_config $name

: ${vboxnet_enable="NO"}

run_rc_command "$1"
