# src/freenas/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/cli
StandardInput=tty
StandardOutput=tty
