# src/freenas/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/cli --menu --pager
StandardInput=tty
StandardOutput=tty
# If someone presses and holds Ctrl-c on the console, we will soon end with
# `getty@tty1.service: Start request repeated too quickly.` and unusable console.
# Let's avoid this by restarting indefinitely.
StartLimitBurst=0
# Do not overload the system with momentary infinite restarts in case something goes wrong
RestartSec=5
