# src/freenas/etc/systemd/system/scst.service.d/override.conf
[Service]
EnvironmentFile=/etc/scst.env
