# src/freenas/etc/systemd/system/uuidd.service.d/override.conf
# we dont create a uuidd user/group so override with nobody/nogroup
[Service]
User=nobody
Group=nogroup
