# src/freenas/etc/ix.rc.d/ix-preinit
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-preinit
# REQUIRE: FILESYSTEMS
# BEFORE: SERVERS

. /etc/rc.subr

do_preinit()
{
	/usr/local/bin/midclt call -job initshutdownscript.execute_init_tasks PREINIT > /dev/null 2>&1
}

name="ix-preinit"
start_cmd='do_preinit'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
