# src/freenas/etc/ix.rc.d/ix-syncdisks
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-syncdisks
# REQUIRE: FILESYSTEMS

. /etc/rc.subr

syncdisks()
{
	touch /tmp/.sync_disk_done
	/usr/local/bin/midclt call notifier.sync_disks
	/usr/local/bin/midclt call notifier.multipath_sync
	checkyesno failover_enable || /usr/local/bin/midclt call notifier.zpool_enclosure_sync 2> /dev/null
}

name="ix-syncdisks"
start_cmd='syncdisks'
stop_cmd=''

load_rc_config $name
run_rc_command "$1"
