# src/freenas/etc/ix.rc.d/ix-sed
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-sed
# REQUIRE: ix-syncdisks
# BEFORE: ix-zfs

. /etc/rc.subr

unlock_drives()
{
	LD_LIBRARY_PATH=/usr/local/lib /usr/local/bin/midclt call disk.sed_unlock_all > /dev/null
}

name="ix-sed"
start_cmd='unlock_drives'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
