# src/freenas/etc/systemd/system/keepalived.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/keepalived --dont-fork --log-detail --vrrp
OOMScoreAdjust=-1000
