# nas_ports/freenas/zfsd/pkg-plist
etc/rc.d/zfsd
etc/zfs/cases
sbin/zfsd
usr/include/devctl/consumer.h
usr/include/devctl/event.h
usr/include/devctl/event_factory.h
usr/include/devctl/exception.h
usr/include/devctl/guid.h
usr/lib/libdevdctl.a
usr/lib/libdevdctl.so
usr/lib/libdevdctl.so.0
usr/lib/libdevdctl_p.a
@dir usr/include/devctl
@dir data/zfs/cases
