# nanobsd/Files/etc/rc.d/ix-afpd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-afpd
# REQUIRE: FILESYSTEMS
# BEFORE: afpd

. /etc/rc.subr

# Defaults
afpd_enable=${afpd_enable:-"NO"}
afpd_config=${afpd_config:-"/usr/local/etc/afpd.conf"}
afpd_avdefault=${afpd_avdefault:-"/usr/local/etc/AppleVolumes.default"}
afpd_avsystem=${afpd_avsystem:-"/usr/local/etc/AppleVolumes.system"}
avahi_daemon_config=${avahi_daemon_config:-"/usr/local/etc/avahi/avahi-daemon.conf"}

# If $1 == 1, return $2, otherwise return nothing
bool_on()
{
    if [ $1 -gt 0 ]; then
	echo -n $2
    fi
}

# If $1 == 1, return $2, otherwise return nothing
bool_off()
{
    if [ $1 -eq 0 ]; then
	echo -n $2
    fi
}

# Return "on" for '1' and "off" for '0'
on_off()
{
    case $1 in
	0) echo "off";;
	1) echo "on";;
    esac
}

generate_afpd_config()
{
    local IFS=\|
    local f="afp_srv_name afp_srv_guest afp_srv_local afp_srv_ddp afp_srv_guest_user"
    eval local $f
    local sf=$(echo $f | sed -e 's/ /, /g')

    ${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} \
        "SELECT $sf FROM services_afp ORDER BY -id LIMIT 1" | \
    while eval read $f; do
	: ${afp_srv_name:=-}
	cat <<EOF
"${afp_srv_name}" $(bool_off ${afp_srv_ddp} "-noddp") -nosavepassword \\
-defaultvol ${afpd_avdefault} -systemvol ${afpd_avsystem} \\
-uservol -uampath /usr/local/libexec/netatalk-uams \\
-guestname ${afp_srv_guest_user} \\
EOF
	if [ ${afp_srv_guest} -eq 1 -a ${afp_srv_local} -eq 0 ]; then
	    echo "-uamlist uams_guest.so"
	elif [ ${afp_srv_local} -eq 1 ]; then
	    echo "-uamlist uams_guest.so,uams_clrtxt.so,uams_dhx.so,uams_randnum.so"
	else
	    echo
	fi
    done
}

generate_av_system()
{
}

generate_av_default()
{
    local IFS=\|
    local f="afp_name afp_comment afp_sharepw afp_sharecharset afp_allow afp_deny afp_ro afp_rw afp_diskdiscovery afp_discoverymode afp_dbpath afp_cachecnid afp_crlf afp_mswindows afp_noadouble afp_nodev afp_nofileid afp_nohex afp_prodos afp_nostat afp_upriv"
    eval local $f
    local sf=$(echo $f | sed -e 's/ /, /g')
    ${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} \
        "SELECT ${sf}, mp_path FROM sharing_afp_share AS aps LEFT OUTER JOIN storage_mountpoint AS mp ON aps.afp_path_id = mp.id ORDER BY aps.id DESC" | \
    while eval read $f mountpoint; do
	if [ -d "${mountpoint}" ]; then
		echo "# ${afp_comment}"
		echo -n "\"${mountpoint}\" \"${afp_name}\""
		echo -n "${afp_sharepw:+ password:${afp_sharepw}}"
#		echo -n "${casefold:+ casefold:${casefold}}"
		echo -n "${afp_sharecharsets:+ sharecharsets:${sharecharsets}}"
		echo -n " options:usedots"
		for _i in cachecnid crlf mswindows noadouble nodev nofileid nohex prodos nostat upriv; do
		    bool_on $(eval echo \$afp_${_i}) ",${_i}"
	        done
		if [ ${afp_diskdiscovery} -eq 1 -a ${afp_discoverymode} = "time-machine" ]; then
		    echo -n ",tm"
		fi
		echo -n "${afp_allow:+ allow:${afp_allow}}"
		echo -n "${afp_deny:+ deny:${afp_deny}}"
		echo -n "${afp_ro:+ rolist:${afp_ro}}"
		echo -n "${afp_rw:+ rwlist:${afp_rw}}"
		echo -n "${afp_dbpath:+ dbpath:\"${afp_dbpath}\"}"
#			-i "string-length(cnidscheme) > 0" -v "concat('cnidscheme:',cnidscheme,' ')" -b \
		echo ""
	fi
    done
}

generate_avahi_daemon_config()
{
    local IFS=\|
    local afp_srv_name
    ${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} \
        "SELECT afp_srv_name FROM services_afp ORDER BY -id LIMIT 1" | \
    while read afp_srv_name; do
	cat << EOF
[server]
host-name=${afp_srv_name}
domain-name=local
use-ipv4=yes
use-ipv6=yes
ratelimit-interval-usec=1000000
ratelimit-burst=1000

[wide-area]
enable-wide-area=yes

[publish]
#disable-publishing=no
#disable-user-service-publishing=no
#add-service-cookie=no
#publish-addresses=yes
#publish-hinfo=yes
#publish-workstation=yes
#publish-domain=yes
#publish-dns-servers=192.168.50.1, 192.168.50.2
#publish-resolv-conf-dns-servers=yes
#publish-aaaa-on-ipv4=yes
#publish-a-on-ipv6=no
EOF
	done
}

generate_afpd()
{
    generate_afpd_config > ${afpd_config}
    generate_av_system > ${afpd_avsystem}
    generate_av_default > ${afpd_avdefault}
    generate_avahi_daemon_config > ${avahi_daemon_config}
}

name="ix-afpd"
start_cmd='generate_afpd'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
