# nanobsd/Files/etc/rc.d/ix-localtime
#!/bin/sh
#
#

# PROVIDE: ix-localtime
# REQUIRE: FILESYSTEMS
# BEFORE: NETWORKING

. /etc/rc.subr

#
# Set up timezone
#
install_timezone()
{
	local IFS="|"
	local f="
		stg_timezone
		"
	eval local $f
	local sf=$(var_to_sf $f)
	${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} "SELECT $sf FROM system_settings ORDER BY -id LIMIT 1" | \
	while read $f; do
		timezone=$(echo ${stg_timezone} | sed -es=america-los_angeles=America/Los_Angeles=)
		if [ -z "${stg_timezone}" ]; then
			stg_timezone=America/Los_Angeles
		fi
		cp /usr/share/zoneinfo/${stg_timezone} /etc/localtime
		echo ${stg_timezone} > /var/db/zoneinfo
	done
}

name="ix-localtime"
start_cmd='install_timezone'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
