# nanobsd/Files/etc/rc.d/ix-nfsd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-nfsd
# REQUIRE: FILESYSTEMS
# BEFORE: nfsd mountd

. /etc/rc.subr

#
# TODO: For now we overwrite /etc/exports, it's desirable to teach mountd about another
# place so user can write their own entry.
#
generate_exports()
{
	local IFS=\|

	local mountpoint allroot network alldirs nfsro quiet
	${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} "SELECT mountpoint, allroot, network, alldirs, nfsro, quiet FROM freenas_unixshare AS us LEFT OUTER JOIN freenas_mountpoint AS mp ON us.unixshare_path_id = mp.volumeid_id ORDER BY us.id DESC" | \
	while read mountpoint allroot network alldirs nfsro quiet; do
		echo -n ${mountpoint}
		if [ -n "${alldirs}" ]; then
			echo -n " -alldirs"
		fi
		if [ -n "${nfsro}" ]; then
			echo -n " -ro"
		fi
		if [ -n "${quiet}" ]; then
			echo -n " -quiet"
		fi
		# Current schema can only do ON and OFF.
		# TODO: Support maproot=uid:gid when schema change is ready.
		if [ "${allroot}" = "1" ]; then
			echo -n " -maproot=0"
		fi
		if [ -n "${network}" ]; then
			echo -n " -network ${network}"
		fi
		echo
	done > /etc/exports
}

name="ix-nfsd"
start_cmd='generate_exports'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
