# nanobsd/Files/etc/rc.d/ix-sysctl
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-sysctl
# REQUIRE: ix-update sysctl

. /etc/rc.subr

set_sysctl()
{
	local IFS="|"
	local f="sysctl_mib sysctl_value"
	eval local $f
	local sf=$(var_to_sf $f)

	# Generate sysctls
	${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} \
		"SELECT $sf FROM system_sysctl ORDER BY id" | \
	while eval read $f; do
		sysctl ${sysctl_mib}="${sysctl_value}"
	done
}

name="ix-sysctl"
start_cmd='set_sysctl'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
