# nanobsd/Files/etc/rc.d/ix-smartd
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-smartd
# REQUIRE: FILESYSTEMS
# BEFORE: smartd

. /etc/rc.subr

: ${smartd_config="/usr/local/etc/smartd.conf"}

generate_smartd_cfg()
{
	local IFS="|"
	local f="
		smart_powermode
		smart_difference
		smart_informal
		smart_critical
		smart_email
	"
	eval local $f
	local sf=$(var_to_sf $f)
	local f2="
		disk_name
		disk_smartoptions
		smarttest_type
		smarttest_hour
		smarttest_daymonth
		smarttest_month
		smarttest_dayweek
		"
	cat <<EOF
################################################
# smartd.conf generated by FreeNAS
################################################
EOF

	local ndaymonth ndayweek nhour nmonth

	${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} "SELECT $sf FROM services_smart ORDER by -id LIMIT 1" |
	while eval read $f; do

		${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} "SELECT $sf2 FROM storage_disk d LEFT JOIN system_smarttest s ON d.id=s.smarttest_disk_id WHERE disk_togglesmart=1 AND disk_enabled=1" |
		while eval read $f2; do
			case "$disk_name" in
			None|*zvol*)
				continue
				;;
			esac
			if [ -z "$(/usr/local/sbin/smartctl -i /dev/${disk_name}|grep "SMART support is: Enabled")" ]; then
				continue
			fi
			echo -n "/dev/${disk_name}"
			echo -n " -n ${smart_powermode}"
			echo -n " -W ${smart_difference},${smart_informal},${smart_critical}"
			if [ -n "${smart_email}" ]; then
				echo -n " -m ${smart_email}"
			fi
			if [ -n "${smarttest_type}" ]; then

			nhour=$(echo "${smarttest_hour}"|awk -F ',' '{ if($0 ~ /^\*\/[0-9]+$/) {
				printf "(";
				num = int(substr($0, 3));
				for(i=0;i<24;i+=num) {
					printf "%.2d", i;
					if(i+num < 24) printf "|"
				}
				printf ")"
			} else if($0 == "*") {
				printf "..";
			} else if($0 != "..") {
				printf "(";
				for(i=1;i<NF+1;i++) {
					printf "%.2d", $i;
					if(i != NF) printf "|"
				};
				printf ")"
			} else printf $0 }')
			ndayweek=$(echo "${smarttest_dayweek}"|tr -s "," "|"|awk '{if($0 != "..") printf "(%s)", $0; else printf $0 }')
			nmonth=$(echo "${smarttest_month}"|awk -F ',' '{ if($0 != "..") {
				printf "(";
				for(i=1;i<NF+1;i++) {
					printf "%.2d", $i;
					if(i != NF) printf "|"
				};
				printf ")"
			} else printf $0 }')
			ndaymonth=$(echo "${smarttest_daymonth}"|awk -F ',' '{ if($0 ~ /^\*\/[0-9]+$/) {
				printf "(";
				num = int(substr($0, 3));
				for(i=1;i<32;i+=num) {
					printf "%.2d", i;
					if(i+num < 32) printf "|"
				}
				printf ")"
			} else if($0 == "*") {
				printf "..";
			} else if($0 != "..") {
				printf "(";
				for(i=1;i<NF+1;i++) {
					printf "%.2d", $i;
					if(i != NF) printf "|"
				};
				printf ")"
			} else printf $0 }')

			echo -n " -s ${smarttest_type}/${nmonth}/${ndaymonth}/${ndayweek}/${nhour}"
			fi
			echo " $disk_smartoptions"
		done
	done
}

generate_smartd()
{
	generate_smartd_cfg > ${smartd_config}
}

name=ix-smartd
start_cmd=generate_smartd
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
