# nanobsd/Files/etc/rc.d/ix-sysctl
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-sysctl
# REQUIRE: ix-update sysctl

. /etc/rc.subr

set_autotune_sysctls()
{
	local IFS="
"
	local autotune
	local tmp

	export PATH=$PATH:/usr/local/bin:/usr/local/sbin

	autotune=/usr/local/bin/autotune
	if [ -x $autotune ]; then
		tmp=$(mktemp /tmp/tmp.XXXXXX)

		$autotune \
			--kernel-reserved=1048576 \
			--userland-reserved=2097152 \
			--conf-file=/etc/sysctl.conf > $tmp

		if [ $? -eq 0 ]; then
			grep -v '^#' $tmp | \
			egrep -v '^[[:space:]]*$' | \
			sed -e "s/'//g" -e 's/"//g' | \
			while read -r line
			do
				[ -n "$line" ] || continue
				sysctl "$line"
			done
		fi
	fi
}

set_user_sysctls()
{
	local IFS="|"
	local f="sysctl_mib sysctl_value"
	eval local $f
	local sf=$(var_to_sf $f)

	# Generate sysctls
	${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} \
		"SELECT $sf FROM system_sysctl WHERE sysctl_enabled = 1 ORDER BY id" | \
	while eval read -r $f; do
		sysctl ${sysctl_mib}="${sysctl_value}"
	done
}

start_sysctl()
{
	set_autotune_sysctls
	set_user_sysctls
}

name="ix-sysctl"
start_cmd='start_sysctl'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
