# nanobsd/Files/etc/rc.d/ix-zfs
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: ix-zfs
# REQUIRE: hostid mountcritlocal
# BEFORE: zfs

. /etc/rc.subr

#
# Generate fstab right before mountlate.
#
import_zpools()
{
	local IFS="|"
	local f="
		vol_name
		vol_guid
	"
	local sf=$(var_to_sf $f)
	if [ ! -d "/data/zfs" ]; then
		mkdir /data/zfs || true
	fi
	${FREENAS_SQLITE_CMD} ${FREENAS_CONFIG} "SELECT $sf FROM storage_volume WHERE vol_fstype = 'ZFS'" | \
	while eval read -r $f; do
		if [ -n "${vol_guid}" ]; then
			/sbin/zpool import -o cachefile=none -R /mnt -f ${vol_guid}
		else
			/sbin/zpool import -o cachefile=none -R /mnt -f ${vol_name}
		fi
		/sbin/zpool set cachefile=/data/zfs/zpool.cache ${vol_name}
		# Fixup mountpoints
		[ -d /mnt/mnt ] && /sbin/zfs inherit -r mountpoint ${vol_name}
	done
}

name="ix-zfs"
start_cmd='import_zpools'
stop_cmd=':'

load_rc_config $name
run_rc_command "$1"
