# src/freenas/etc/systemd/system/nfs-ganesha.service.d/override.conf
[Service]
ExecStartPre=-mkdir -p /var/log/ganesha
