# nas_ports/shells/scponly/pkg-plist
@shell bin/scponly
%%CHROOT%%@exec echo ""
%%CHROOT%%@exec echo "To setup chroot cage, run the following commands:"
%%CHROOT%%@exec echo "  1) cd %%PREFIX%%/%%EXAMPLESDIR%%/ && /bin/sh setup_chroot.sh"
%%CHROOT%%@exec echo "  2) Set scponlyc_enable=\"YES\" in /etc/rc.conf"
%%CHROOT%%@exec echo "  3) Run %%PREFIX%%/etc/rc.d/scponly start"
%%PORTDOCS%%@exec echo ""
%%PORTDOCS%%@exec echo "For information on several potential security concerns,"
%%PORTDOCS%%@exec echo "please read:"
%%PORTDOCS%%@exec echo "%%PREFIX%%/%%DOCSDIR%%/SECURITY"
%%CHROOT%%@shell sbin/scponlyc
%%CHROOT%%%%EXAMPLESDIR%%/setup_chroot.sh
%%CHROOT%%%%EXAMPLESDIR%%/config.h
etc/scponly/debuglevel
man/man8/scponly.8.gz
